# docker_root/etc/services.d/cyclops/run
#!/usr/bin/execlineb -P
/usr/bin/python3 /usr/local/src/cyclops/core.py
